# ansible/roles/jupyter/files/jupyterhub.service
[Unit]
Description=Jupyterhub

[Service]
ExecStart=/usr/local/bin/jupyterhub
WorkingDirectory=/etc/jupyterhub

[Install]
WantedBy=multi-user.target
